# systemd/monitor_agent.service
[Unit]
Description=EC2 Monitoring Agent
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/monitor_agent
ExecStart=/opt/monitor_agent/venv/bin/python agent.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
